# src/nexus/service/installation/nexus_service.service


[Unit]
Description=Nexus GPU Manager
After=network.target

[Service]
Type=simple
User=nexus
Group=nexus
WorkingDirectory=/home/nexus
# Adjust the ExecStart path to point to the installed nexus-service executable
ExecStart=/usr/local/bin/nexus-service
Restart=on-failure
Environment=PYTHONUNBUFFERED=1
# Optionally, include an EnvironmentFile:
# EnvironmentFile=/etc/default/nexus

[Install]
WantedBy=multi-user.target
